# nix/ci/systemd/unit-files/default.target
[Unit]
Description=Default target
Requires=@sysinit_target_name@ @systemd_journald_service_name@
