# nix/ci/systemd/unit-files/halt.service
[Unit]
Description=Halt service
DefaultDependencies=no

[Service]
ExecStart=@systemd@/bin/systemctl --force halt
