# nix/ci/systemd/unit-files/halt.target
[Unit]
Description=Halt target
Requires=@halt_service_name@
