# nix/ci/systemd/unit-files/systemd-journald.socket
[Unit]
Description=Journal socket
DefaultDependencies=no

[Socket]
ListenStream=/run/systemd/journal/stdout
ListenDatagram=/run/systemd/journal/socket
